# cookbooks/task1_jboss/templates/jboss.service.erb
[Unit]
Description=JBoss Service
After=network.target

[Service]
Type=idle
ExecStart=<%= @entrypoint %>
RemainAfterExit=yes
Environment="JAVA_HOME=<%= @java_home %>"
Environment="JBOSS_HOME=<%= @jboss_home %>"
TimeoutStartSec=600
TimeoutStopSec=600

[Install]
WantedBy=multi-user.target
